# autobuild/unified/filogic/mac80211/25.12/files/target/linux/mediatek/base-files/etc/uci-defaults/99-enable-wifi-teardown
chmod +x /etc/init.d/wifi-teardown
/etc/init.d/wifi-teardown enable

exit 0
