# autobuild/unified/filogic/mac80211/25.12/files/target/linux/mediatek/base-files/etc/init.d/wifi-teardown
#!/bin/sh /etc/rc.common

START=99
STOP=01

start() {
    return 0
}

stop() {
    echo "wifi-teardown: executing wifi down" > /dev/kmsg

    wifi down
    sleep 1
}
